# systemd/.config/systemd/user/backup_ehforwarderbot_wechat.timer
[Unit]
Description=Run backup_ehforwarderbot weekly

[Timer]
OnBootSec=33min
OnUnitActiveSec=1w 
 
[Install]
WantedBy=timers.target
