# systemd/.config/systemd/user/delete_vim_tags.timer
[Unit]
Description=Run delete_vim_tags weekly

[Timer]
OnBootSec=43min
OnUnitActiveSec=3M 
 
[Install]
WantedBy=timers.target
