# systemd/.config/systemd/user/delete_vim_undo_dir.timer
[Unit]
Description=Run delete_vim_undo_dir weekly

[Timer]
OnBootSec=47min
OnUnitActiveSec=7w 
 
[Install]
WantedBy=timers.target
